# guide/service/service-proxy/run
#!/bin/sh

exec docker run \
    --net=host \
    --privileged \
    pwittrock/hyperkube \
      /hyperkube proxy \
      --master=http://127.0.0.1:8080 \
      --v=0